# sophia-bot.service
[Unit]
Description=Sophia Telegram Bot
After=network.target

[Service]
Type=simple
# Edit these to match your setup:
User=ubuntu
WorkingDirectory=/home/ubuntu/sophia
Environment=PATH=/home/ubuntu/.local/bin:/home/ubuntu/.cargo/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ubuntu/sophia/target/release/sophia
# Restart always: on failure AND on exit code 42 (auto-update)
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
